# debian/postinst
#!/bin/sh
# postinst script for lwsyslogger

case "$1" in
    configure)
      /bin/systemctl daemon-reload
      /usr/sbin/groupadd -r -f lwsyslogger
      /usr/sbin/useradd -r -g lwsyslogger -s /bin/false -c "lwsyslogger system" lwsyslogger
      if test ! -e /var/log/lwsyslogger ; then
	  mkdir -p /var/log/lwsyslogger
	  chown lwsyslogger:lwsyslogger /var/log/lwsyslogger
	  chmod 770 /var/log/lwsyslogger
      fi
      if test ! -d /etc/lwsyslogger.d ; then
	 mkdir -p /etc/lwsyslogger.d
      fi
      if test ! -e /etc/lwsyslogger.conf ; then
	  cp /usr/share/lwsyslogger/lwsyslogger-sample.conf /etc/lwsyslogger.conf
      fi
      /bin/systemctl restart lwsyslogger
      /bin/systemctl enable lwsyslogger
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
